# crontab
* * * * * cd /home/scraper && python src/scraper.py
